# usr/lib/systemd/system/iptv.timer
[Unit]
Description=Runs iptv every hour

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
Unit=iptv.service

[Install]
WantedBy=timers.target
