# usr/lib/systemd/system/iptv.service
[Unit]
Description=Create iptv m3u with all multiplexes

[Service]
Type=simple
ExecStart=/etc/mumudvb/iptv.sh

[Install]
WantedBy=multi-user.target
